# deploy/systemd/family-summarizer.timer
[Unit]
Description=Family inc. WhatsApp summarizer timer (hourly, 24h)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
