# deploy/systemd/family-weekly.timer
[Unit]
Description=Family inc. weekly briefing timer (Sat 21:00 Asia/Jerusalem)

[Timer]
OnCalendar=Sat *-*-* 21:00:00
Persistent=true

[Install]
WantedBy=timers.target
